# configs/safehaven-wlan1-ip.service
[Unit]
# SafeHaven Pi — Auto-recover wlan1 gateway IP (10.42.0.1/24) on device
# appearance. Handles both cold boot AND hot-plug cycles (USB adapter
# knocked loose, reseated, etc).
#
# Lifecycle is tied directly to the wlan1 network-device systemd unit:
#
#   BindsTo   → if wlan1 disappears, we stop
#   PartOf    → restart/stop propagates through from the device unit
#   WantedBy  → wlan1 device appearing pulls us in
#   After     → we run AFTER wlan1 exists (not racing with it)
#
# After assigning the IP we also kick hostapd + dnsmasq so they grab the
# fresh interface cleanly. Without this, hostapd can hold a stale
# reference to the old (now-gone) wlan1 and silently fail to broadcast.
#
# Installed to /etc/systemd/system/ by install.sh.

Description=SafeHaven — auto-recover wlan1 gateway IP and hotspot
BindsTo=sys-subsystem-net-devices-wlan1.device
PartOf=sys-subsystem-net-devices-wlan1.device
After=sys-subsystem-net-devices-wlan1.device
After=hostapd.service
Wants=hostapd.service

[Service]
Type=oneshot
# Brief pause so the device settles (USB re-enumeration, driver bind, etc)
ExecStartPre=/bin/sleep 2
# Flush any stale addresses, then assign the gateway. Dash prefix = ignore
# error if it's already set (idempotent for multi-trigger cases).
ExecStart=-/usr/sbin/ip addr flush dev wlan1
ExecStart=/usr/sbin/ip addr add 10.42.0.1/24 dev wlan1
# Kick hostapd and dnsmasq so they re-attach to the fresh interface.
# try-restart = only restart if already running, otherwise no-op.
ExecStartPost=-/bin/systemctl try-restart hostapd.service
ExecStartPost=-/bin/systemctl try-restart dnsmasq.service
# Clean up on stop (device removed or service stopped)
ExecStop=-/usr/sbin/ip addr del 10.42.0.1/24 dev wlan1
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
# Also pulled in whenever wlan1 appears, so hot-plug re-triggers us
WantedBy=sys-subsystem-net-devices-wlan1.device
